# WM/hyprland/autostart/sway-notification-center.desktop
[Desktop Entry]
Type=Application
Version=1.5
Name=Sway Notification Center
Exec=swaync
Terminal=false
