# WM/hyprland/autostart/swww-daemon.desktop
[Desktop Entry]
Type=Application
Version=1.5
Name=swww-daemon
GenericName=Wallpaper wrapper
Exec=swww-daemon
Terminal=false
