# WM/hyprland/autostart/udiskie.desktop
[Desktop Entry]
Type=Application
Version=1.5
Name=udiskie
Exec=udiskie
Comment[ru]=Автомонтирование usb накопителей
Terminal=false
