# WM/hyprland/autostart/waybar.desktop
[Desktop Entry]
Type=Application
Version=1.5
Name=Waybar
GenericName=Bar
Exec=waybar
Terminal=false
